# supervisor-sshchrooted.conf
[program:sshchrooted]
command=/root/scripts/ssh-chrooted-setup.sh
numprocs=1
autostart=true
autorestart=false
startsecs=0